# config/memory_watch.service
[Unit]
Description=memory watch script

[Service]
ExecStart=/home/ubuntu/cellxgene_ec2/src/memory_watch.sh

[Install]
WantedBy=multi-user.target
